# docs/supervisor.conf
; Supervisor config - ChatBot Blasternet workers
; Copie este arquivo para /etc/supervisor/conf.d/chatbot-workers.conf
; e ajuste "directory", "command" e "user" conforme seu servidor.

[program:chatbot-worker-default]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/ChatBot-Blasternet/backend
command=php artisan queue:work redis --queue=default,emails,webhooks,campaigns --sleep=3 --tries=3 --timeout=180 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chatbot-worker-default.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=3600

[program:chatbot-worker-realtime]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/ChatBot-Blasternet/backend
command=php artisan queue:work redis --queue=realtime --sleep=1 --tries=3 --timeout=120 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chatbot-worker-realtime.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=3600

[program:chatbot-worker-ai]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/ChatBot-Blasternet/backend
command=php artisan queue:work redis --queue=ai --sleep=1 --tries=2 --timeout=300 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chatbot-worker-ai.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stopwaitsecs=3600

[group:chatbot-workers]
programs=chatbot-worker-default,chatbot-worker-realtime,chatbot-worker-ai
